# contrib/init-scripts/batwatch.openrc.in
#!/sbin/runscript
#
# @@NAME@@ - daemon start/stop script (openrc)
#
# Copyright (C) 2014 Andre Erdmann <dywi@mailerd.de>
# Distributed under the terms of the GNU General Public License;
# either version 2 of the License, or (at your option) any later version.
#

extra_started_commands="reload"
description_reload="forcefully update batteries"

@@VARS@@

depend() {
   use logger
   need dbus dev localmount
   after netmount nfsmount
}

@@F_EXE_CHECK@@

batwatch_prepare() {
   if ! @@EXE_CHECK@@; then
      eerror "At least one --exe must be specified in \$BATWATCH_OPTS"
      return 4
   fi

   case "${BATWATCH_PIDFILE}" in
      */batwatch/?*)
         checkpath -d -m 0750 \
            -o "${BATWATCH_USER:-root}:${BATWATCH_GROUP:-root}" \
            "${BATWATCH_PIDFILE%/*}"
      ;;
   esac
}

@@STOP@@() {
   @@STOP_CODE@@
}

start() {
   case "${RC_CMD}" in
      reload)
         true
      ;;
      *)
         batwatch_prepare || return 1
      ;;
   esac
   ebegin "Starting ${RC_SVCNAME}"
   @@START_CODE@@
   eend ${?}
}

stop() {
   ebegin "Stopping ${RC_SVCNAME}"
   @@STOP@@
   eend ${?}
}

reload() {
   batwatch_prepare || return 1
   ebegin "Reloading ${RC_SVNAME}"
   @@STOP@@ @@RELOAD_ARGS@@
   eend ${?}
}
